# linux/debian/mozillavpn.service
[Unit]
Description=MozillaVPN D-Bus service

[Service]
Type=dbus
BusName=org.mozilla.vpn.dbus
ExecStart=/usr/bin/mozillavpn linuxdaemon
