# linux/mozillavpn.service.in
[Unit]
Description=MozillaVPN D-Bus service
After=modprobe@xt_cgroup.service
Requires=modprobe@xt_cgroup.service

[Service]
Type=dbus
BusName=org.mozilla.vpn.dbus
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/mozillavpn linuxdaemon

[Install]
WantedBy=multi-user.target
