# linux/mozillavpn.service.in
[Unit]
Description=MozillaVPN D-Bus service
Wants=modprobe@xt_cgroup.service
After=modprobe@xt_cgroup.service

[Service]
Type=dbus
BusName=org.mozilla.vpn.dbus
Restart=on-failure
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/mozillavpn linuxdaemon

[Install]
WantedBy=multi-user.target
