# extension/socks5proxy/bin/socksproxy.service.in
[Unit]
Description=MozillaVPN SOCKS Proxy service

[Service]
Type=simple
Restart=on-failure
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/socksproxy --unix /var/run/mozillavpn.proxy
ExecStopPost=rm -f /var/run/mozillavpn.proxy

[Install]
WantedBy=multi-user.target
